# main.cron
@reboot /bin/bash /u01/aipoc/No.1-RAG/main.sh &