# main.cron
@reboot /bin/bash /home/oracle/No.1-RAG/main.sh